# scripts/cron/crontab
0 0* * * root bash /app/scripts/cron/delete_old_threads.sh
